# deploy/raspberry-pilot-kiosk/kiosk-master.service
[Unit]
Description=KOR35 Kiosk Master (dual HDMI + Chromium)
After=display-manager.service graphical.target network.target
# Non usare network-online.target: in bosco senza WAN il boot resta bloccato.

[Service]
Type=simple
User=pi
Group=pi
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/pi/.Xauthority
Environment=DBUS_SESSION_BUS_ADDRESS=unix:path=/run/user/1000/bus
EnvironmentFile=-/etc/kor35/kiosk.env
ExecStartPre=/bin/sleep 6
ExecStart=/usr/local/bin/kiosk-master.sh
Restart=on-failure
RestartSec=12

[Install]
WantedBy=graphical.target
