# deploy/raspberry-pilot-kiosk/kor35-kiosk.service
[Unit]
Description=KOR35 Pilot Console Kiosk (client verso Pi mirror)
# Pi dedicato ai display: NON dipende da docker/mirror-stack (sono su un altro host).
# Evitare network-online.target: offline in evento il boot resta bloccato senza uplink.
After=graphical.target network.target

[Service]
Type=simple
User=pi
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/pi/.Xauthority
EnvironmentFile=-/etc/kor35/kiosk.env
ExecStart=/usr/local/bin/kiosk-master.sh
Restart=always
RestartSec=5

[Install]
WantedBy=graphical.target

# Deprecato: usare kiosk-master.service. Mantenuto per compatibilità installazioni vecchie.
