# deploy/raspberry-pilot-kiosk/kor35-kiosk.service
[Unit]
Description=KOR35 Pilot Console Kiosk (client verso Pi mirror)
# Pi dedicato ai display: NON dipende da docker/mirror-stack (sono su un altro host).
# Evitare network-online.target: offline in evento il boot resta bloccato senza uplink.
After=graphical.target network.target

[Service]
Type=simple
User=pi
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/pi/.Xauthority
# IP del Pi mirror (gateway WiFi evento / stack Docker). Adattare se diverso da 10.42.0.1
Environment=PILOT_BASE_URL=http://10.42.0.1
EnvironmentFile=-/etc/kor35/kiosk.env
ExecStartPre=/bin/bash -lc 'base="${PILOT_BASE_URL%/}"; for _ in $(seq 1 150); do curl -fsS "${base}/api/healthz/" >/dev/null && exit 0; sleep 2; done; echo "Mirror non raggiungibile: ${base}/api/healthz/" >&2; exit 1'
ExecStart=/opt/kor35-kiosk/start-kiosk.sh
Restart=always
RestartSec=5

[Install]
WantedBy=graphical.target
