# config/systemd/kor35-mirror-resync.service
[Unit]
Description=KOR35 Mirror - Riallineamento post-evento (DB+Media)
After=network.target docker.service kor35-mirror-stack.service

[Service]
Type=oneshot
User=pi
Group=pi

# Adatta se il path repo sul Pi e' diverso
WorkingDirectory=/home/pi/kor35-replica
Environment=COMPOSE_PROJECT_NAME=kor35-replica
Environment=KOR35_BACKEND_ENV_FILE=/home/pi/kor35-replica/backend/.env.mirror

# Richiede .env.sync-media configurato per l'accesso SSH al master.
ExecStart=/home/pi/kor35-replica/scripts/mirror_resync_after_event.sh --env mirror

NoNewPrivileges=true
PrivateTmp=true
